# scrapers/crontab
0 */4 * * * python3 -m scrapy runspider /app/book_scraper.py
